# projects/start-wrt/openwrt-overlay/target/linux/spacemit/base-files/etc/uci-defaults/98-custom-init-env.sh
#!/bin/sh

echo "custom init env"

ln -s /lib/libc.so /usr/lib/libc.so

exit 0
