# projects/start-wrt/openwrt-overlay/target/linux/spacemit/base-files/etc/uci-defaults/99-custom-ca-certificates
#!/bin/sh
cat /etc/ssl/certs/dc.com-CA-Root-Base64.crt >> /etc/ssl/certs/ca-certificates.crt
exit 0

